# 13.service
[Unit]
Description=13

[Service]
Type=simple
ExecStart=/var/www/html/not/13.php

[Install]
WantedBy=multi-user.target
